# rfmp-daemon/rfmpd.service
[Unit]
Description=RFMP Daemon - RF Microblog Protocol
After=network.target direwolf.service
Wants=direwolf.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/IdeaProjects/rfmp-daemon
Environment="PATH=/home/pi/IdeaProjects/rfmp-daemon/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Create runtime directories
ExecStartPre=/bin/mkdir -p /home/pi/rfmpd
ExecStartPre=/bin/chown pi:pi /home/pi/rfmpd

# Start the daemon
ExecStart=/home/pi/IdeaProjects/rfmp-daemon/venv/bin/python -m rfmpd.main -c /home/pi/rfmpd/config.yaml

# Restart on failure
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target